# Chapter18/Exercise18.05/bookr/service-conf/gunicorn-bookr.socket
[Unit]
Description=Bookr Gunicorn Socket

[Socket]
ListenStream=/run/gunicorn-bookr.sock
User=www-data
Mode=600

[Install]
WantedBy=sockets.target
